# card-server.service
[Unit]
Description=ShortFactory Card Server (Node.js + LLaMA)
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/vhosts/shortfactory.shop/httpdocs
ExecStart=/usr/bin/node /var/www/vhosts/shortfactory.shop/httpdocs/card-server.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
